# filesystem/etc/service/httpd/run
#!/bin/sh
#exec httpd -D FOREGROUND
exec apache2-foreground
