# filesystem/etc/service/httpd/run
#!/bin/sh
exec apache2 -D FOREGROUND
#exec apache2-foreground
